# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:mariadb]
command=mysqld --user=mysql
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/mysql/mariadb.log
stderr_logfile=/var/log/mysql/mariadb_err.log
startsecs=5
stopwaitsecs=30
stopasgroup=true
killasgroup=true
redirect_stderr=true

[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/apache2/supervisor_access.log
stderr_logfile=/var/log/apache2/supervisor_error.log
startsecs=3
stopwaitsecs=5
stopasgroup=true
killasgroup=true
redirect_stderr=true
